# crontab
@hourly ~/dotfiles/bin/remind_agenda > ~/dotfiles/log/remind_agenda.log 2>&1
